# dists/bookworm/dynamic/var/centauri/configs/bookworm/@admin/systemd/system/centauricron-midnight.service
[Unit]
Description=Run centauricron at midnight

[Service]
Type=oneshot
ExecStart=/usr/local/bin/centauricron --append midnight

