# dists/bookworm/dynamic/var/centauri/configs/bookworm/@admin/systemd/system/centauricron-boot.timer
[Unit]
Description=Run centauricron at boot time

[Timer]
OnBootSec=30

[Install]
WantedBy=timers.target
