# dists/bookworm/dynamic/var/centauri/configs/bookworm/@admin/systemd/system/centauricron-hourly.service
[Unit]
Description=Run centauricron every hour

[Service]
Type=oneshot
ExecStart=/usr/local/bin/centauricron --silent --append hourly
