# dists/bookworm/dynamic/var/centauri/configs/bookworm/@admin/systemd/system/centauricron-boot.service
[Unit]
Description=Run centauricron at boot time

[Service]
Type=oneshot
ExecStart=/usr/local/bin/centauricron --append boot

