# dists/bookworm/dynamic/var/centauri/configs/bookworm/@admin/systemd/system/centauricron-midnight.timer
[Unit]
Description=Run centauricron at midnight

[Timer]
OnCalendar=*-*-* 0:00:00
RandomizedDelaySec=90

[Install]
WantedBy=timers.target
